# templates/casparcgserver.service.j2
[Service]
Description=CasparCG Server
ExecStart={{ casparcg_install_path }}/casparcg_server/run.sh
Restart=on-failure
RestartSec=10
