# templates/casparcgserver.service.j2
[Service]
Description=CasparCG Server
ExecStart=/usr/bin/run.sh
Restart=on-failure
RestartSec=10
